# pybtsp.service.template
[Unit]
Description=Python BTSP Server
After=sys-bluetooth.service

[Service]
Type=simple
ExecStart={pwd}/server.py

[Install]
WantedBy=bluetooth.target
